# deploy/aws/systemd/journal.service
[Unit]
Description=Triple Screen Journal FastAPI Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/triple_screen
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/ubuntu/triple_screen/src
ExecStart=/home/ubuntu/triple_screen/.venv/bin/python -m journal
Restart=always
RestartSec=3
StandardOutput=append:/home/ubuntu/triple_screen/logs/journal.log
StandardError=append:/home/ubuntu/triple_screen/logs/journal.log

[Install]
WantedBy=multi-user.target
