# deploy/aws/systemd/triple-screen.service
[Unit]
Description=Triple Screen Scanner One-Shot Run
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ec2-user
WorkingDirectory=/home/ec2-user/triple_screen
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/ec2-user/triple_screen/.venv/bin/python src/scanner.py --once
StandardOutput=append:/home/ec2-user/triple_screen/logs/systemd.log
StandardError=append:/home/ec2-user/triple_screen/logs/systemd.log

[Install]
WantedBy=multi-user.target
